# config/supervisor/mio-ressources-worker.conf
[program:mio-ressources-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/mio-ressources-site/artisan queue:work redis --sleep=3 --tries=3 --max-time=3600
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=2
redirect_stderr=true
stdout_logfile=/var/www/mio-ressources-site/storage/logs/worker.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stopwaitsecs=3600
